# pulse.service
[Unit]
Description=Pulse - Alpha's scheduled job runner
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/Pondside/Basement/Pulse
ExecStart=uv run python -m pulse.main
Restart=always
RestartSec=10

# Logging to journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pulse

[Install]
WantedBy=default.target
