# pulse.service
[Unit]
Description=Pulse - Alpha's scheduled job runner
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/Pondside/Basement/Pulse
# UV_PATH for subprocess spawning (solitude_next, etc.) since systemd doesn't inherit user PATH
Environment="UV_PATH=/home/alpha/.local/bin/uv"
ExecStart=/home/alpha/.local/bin/uv run python -m pulse.main
Restart=always
RestartSec=10

# Logging to journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pulse

[Install]
WantedBy=default.target
